# proxy-nsd@.service.in
[Unit]
Description=proxy-ns
After=network.target

[Service]
ExecStart=@PREFIX@/bin/proxy-nsd
EnvironmentFile=/etc/proxy-nsd/%i.conf

[Install]
WantedBy=multi-user.target
