# proxy-nsd.service
[Unit]
Description=proxy-ns
After=network-online.target

[Service]
ExecStart=/usr/bin/stdbuf -oL /usr/bin/proxy-nsd
EnvironmentFile=/etc/proxy-nsd.conf

[Install]
WantedBy=multi-user.target
